# templates/files-horizon.service
[Unit]
Description=Horizon worker for Raki

[Service]
Type=simple

ExecStart=PATHTOPHP PATHTORAKI/artisan horizon

[Install]
WantedBy=default.target
